# debian/cloudprint-service.cloudprintd.service
[Unit]
Description=Google Cloud Print proxy service
Wants=network-online.target
After=cups.service network.target network-online.target

[Service]
ExecStart=/usr/sbin/cloudprintd -a /var/lib/cloudprintd/authfile.json
Type=simple
Restart=on-failure
RestartSec=300

[Install]
WantedBy=multi-user.target
